# fake-service.ini
[program:consul]
command=/usr/local/bin/consul agent -datacenter %(ENV_CONSUL_DATACENTER)s -config-dir /config -client 0.0.0.0 -retry-join %(ENV_CONSUL_SERVER)s -bind 0.0.0.0 -grpc-port=8502 -data-dir /tmp -hcl "enable_central_service_config = true"
autorestart=true
priority=100

[program:envoy]
command=/usr/local/bin/consul connect envoy -sidecar-for %(ENV_SERVICE_ID)s -envoy-binary /usr/local/bin/envoy -- -l debug
autorestart=true
startretries=300
priority=200

[program:fake-service]
command=/app/fake-service
autorestart=true
startretries=300
priority=300

[program:register-service]
command=/entrypoint.sh
autorestart=false
priority=400
